# root/etc/services.d/noxbit/run
#!/usr/bin/with-contenv bash
exec \
	s6-setuidgid abc ./STM-Hypervisor \
	-config=/config/noxbit.cfg
